# rc.local
#!/bin/sh -e

su - pi -c '/user/bin/vncserver :1 -geometry 1920x1080 -depth 24 -dpi 96'

exit 0

