# rc.local
#!/bin/sh -e

su - pi -c '/usr/bin/tightvncserver :1'

exit 0

